# docker/rootfs/etc/services.d/emulated-hue/finish
#!/usr/bin/execlineb -S0
# ==============================================================================
# Community Hass.io Add-ons: Emulated Hue
# Take down the S6 supervision tree when Emulated Hue fails to start/crashes
# ==============================================================================

s6-svscanctl -t /var/run/s6/services
