# docker/rootfs/etc/services.d/emulated-hue/run
#!/command/with-contenv bashio
# ==============================================================================
# Emulated Hue
# Runs Emulated Hue
# ==============================================================================

bashio::log.info "Starting Emulated Hue..."
cd /app
exec python3 -m emulated_hue
